# roles/docker/files/docker_cleanup.service
[Unit]
Description=Run 'docker system prune -f'

[Service]
ExecStart=/usr/bin/docker system prune -f
