# roles/docker/files/docker_cleanup.timer
[Unit]
Description=Run docker_cleanup.service once a month

[Timer]
OnCalendar=monthly

[Install]
WantedBy=timers.target
