# roles/nfs/templates/krb_mkhomedir.timer.j2
[Unit]
Description=Run krb_mkhomedir.service every 5 minutes

[Timer]
OnBootSec=5min
OnUnitActiveSec=5min

[INstall]
WantedBy=timers.target
