# roles/nfs/templates/krb_mkhomedir.service.j2
[Unit]
Description=Run krb_mkhomedir.sh script

[Service]
ExecStart=/usr/local/bin/krb_mkhomedir.sh
