# roles/nfs/templates/scratch_deletion.service.j2
[Unit]
Description=Delete scratch contents once a month

[Service]
ExecStart=/bin/bash rm -r /nfs/tank/scratch/*
